# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=5MB
logfile_backups=3

[program:medusa]
command=/app/backend/start.sh
directory=/app/backend
autostart=true
autorestart=true
startsecs=30
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/supervisor/medusa.out.log
stderr_logfile=/var/log/supervisor/medusa.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
redirect_stderr=true
environment=NODE_ENV="production",PORT="9000"
priority=10

[program:nextjs]
command=node server.js
directory=/app/storefront
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/supervisor/nextjs.out.log
stderr_logfile=/var/log/supervisor/nextjs.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=NODE_ENV="production",PORT="8000",HOSTNAME="0.0.0.0"
priority=20

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=5
stopwaitsecs=10
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/nginx.out.log
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=30
